# resel-nsa-0.4/debian/resel-nsa.cron.d
# Cron jobs for NSA

* * * * * nsa /usr/libexec/nsa/nsa-collect
* * * * * nsa /usr/libexec/nsa/nsa-send
